# service/led-tailscale-init.service
[Unit]
Description=Ensure the Pi is joined to the tailnet (idempotent)
# `time-sync.target` is what we *really* depend on: tailscale's apt
# repo signatures (and Tailscale auth-key validation later) reject
# requests when the system clock is too far off, and the Pi Zero W
# has no RTC — first boot has wall clock pinned to the raspbian
# image build date until NTP syncs. Without this dep, apt-get
# update bails with "Not live until <future>" on signed packages
# and the service exits 100. See 2026-04-28 led-delta forensics.
After=led-wifi-setup.service network-online.target time-sync.target
Wants=network-online.target time-sync.target
Requires=led-wifi-setup.service
# Note: no ConditionPathExists. Earlier versions of this unit gated
# on /var/lib/led-wifi-setup/tailscale-joined and refused to run
# again after first boot, leaving no recovery path when tailscaled
# flipped to NeedsLogin (e.g. controlplane logged the device out).
# The script is idempotent now — safe to run on every boot.

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/led/init.env
ExecStart=/usr/local/bin/led-tailscale-init

[Install]
WantedBy=multi-user.target
