# service/led-driver.service
[Unit]
Description=LED matrix driver
# No `After=network-online.target` or `After=led-tailscale-init`:
# the driver needs to start as early as possible so it can render
# the boot frame while wifi-setup runs the AP (and then the setup
# frame once /run/led-wifi-setup.active appears). The driver
# tolerates a missing network — the realtime listener auto-reconnects
# and OTel exports time out gracefully.
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/led-driver --config /usr/local/etc/led/config.toml
Restart=always
RestartSec=2s
User=root

# rpi-led-panel drives HUB75 via tight CPU + DMA loops, so preemption
# mid-frame shows up as visible flicker. We *want* SCHED_FIFO here,
# but on the Zero W's single core, FIFO=99 with no budget cap starved
# PID 1 systemd to the point that it couldn't pet the BCM2835 hardware
# watchdog (1-min timeout) and the kernel rebooted before multi-user
# was ever reached — see the led-delta forensics on 2026-04-28 (20
# boots in 100 minutes, none reached multi-user.target). Until we wire
# up CPUQuota + LimitRTPRIO + a tuned kernel.sched_rt_runtime_us, we
# stick with nice/io preference only. Visible flicker may regress on
# the Pi Zero W; the fix path is to add bounded RT later, not to
# reinstate unbounded FIFO.
Nice=-20
IOSchedulingClass=realtime
IOSchedulingPriority=0

[Install]
WantedBy=multi-user.target
